# extras/udev_rules/99-LeapFrog-Net.rules
SUBSYSTEM=="net", ACTION=="add", ATTRS{idVendor}=="0f63", ATTRS{idProduct}=="001[01]", RUN+="/usr/sbin/avahi-autoipd --daemon --force-bind %k"

